# deploy/polymarket-settler.service
[Unit]
Description=Polymarket Bot - Settlement Monitor
After=network.target

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-bot
EnvironmentFile=/opt/polymarket-bot/.env
ExecStart=/opt/polymarket-bot/venv/bin/python run.py --settle --interval=3600
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
